# packaging/scripts/postinstall
#!/bin/sh

systemctl daemon-reload
systemctl enable systemd-resolved-exporter.service
systemctl restart systemd-resolved-exporter.service