# packaging/scripts/preremove
#!/bin/sh

systemctl stop systemd-resolved-exporter.service
systemctl disable systemd-resolved-exporter.service